# overlay/etc/supervisor/conf.d/filebeat.conf
[program:filebeat]
command=/usr/bin/filebeat
stdout_events_enabled=true
stderr_events_enabled=true